# bin/cortex-job-server
#!/bin/sh

# CorTeX Job Server, Linux service
# Copyright (C) 2013 Deyan Ginev
# All rights reserved.
#
# Use and distribution licensed under the BSD license.

### BEGIN INIT INFO
# Provides:          cortex-job-server
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Boot up clients for all registered CorTeX services
### END INIT INFO

prefix=/usr
exec_prefix=${prefix}
NAME=cortex-gears
DAEMON=${exec_prefix}/local/bin/cortex-gears
PIDDIR=/var/run/cortex-job-server
PIDFILE=${PIDDIR}/cortex-job-server.pid
#CORTEXUSER="cortex"
#__CORTEX_DB_DIR__

test -x ${DAEMON} || exit 0

. /lib/lsb/init-functions

start()
{
  log_daemon_msg "Starting CorTeX Job Server" "cortex-gears"
  if ! test -d ${PIDDIR}
  then
    mkdir ${PIDDIR}
    #chown ${CORTEXUSER} ${PIDDIR}
  fi
  if ! test -d ${CORTEX_DB_DIR}
  then
    log_end_msg 1
    log_warning_msg "Please set CORTEX_DB_DIR in cortex-job-server and rerun the make installation process."
    exit 1
  fi
  if start-stop-daemon \
    --start \
    --exec $DAEMON \
    --pidfile=$PIDFILE \
    --make-pidfile \
    --background \
    -- ${CORTEX_DB_DIR}
  then
    log_end_msg 0
  else
    log_end_msg 1
    log_warning_msg "Could not initialize cortex-gears!"
    exit 1
  fi
}

stop()
{
  log_daemon_msg "Stopping CorTeX Job Server" "cortex-gears"
  if start-stop-daemon \
    --stop \
    --pidfile $PIDFILE
  then
    log_end_msg 0
  else
    log_end_msg 1
    exit 1
  fi
}

status()
{
    status_of_proc -p $PIDFILE $DAEMON $NAME && exit 0 || exit $?
}

case "$1" in

  start)
    start
  ;;

  stop)
    stop
  ;;

  status)
    status
  ;;

  restart|force-reload)
    stop
    start
  ;;

  *)
    echo "Usage: $0 {start|stop|restart|force-reload|status}"
  ;;

esac
